# .service
[Unit]
Description=Gunicorn instance to serve Flask-web
After=network.target

[Service]
User=azureuser
Group=www-data
WorkingDirectory=/home/azureuser/Flask-web
Environment="PATH=/home/azureuser/Flask-web/env/bin"
ExecStart=/home/azureuser/Flask-web/env/bin/gunicorn --bind 0.0.0.0:5000 wsgi:app

[Install]
WantedBy=multi-user.target
